# debian/cobbler-web.postinst
#!/bin/sh

if [ "$1" = "configure" ] ; then
    # Need to restart apache to pickup web configs
    invoke-rc.d apache2 restart
fi

#DEBHELPER#
