# debian/cobbler.upstart
# cobbler - provisioning service
#
# Cobbler is a provisioning service

description	"Cobbler"
author "Dustin Kirkland <kirkland@canonical.com>"

start on filesystem and net-device-up
stop on runlevel [016]

respawn

pre-start script
  [ -d "/var/lib/cobbler" ] || mkdir -p /var/lib/cobbler
end script

# To add options to your daemon, edit the line below:
exec /usr/bin/cobblerd --no-daemonize
